# templates/default/bootapp.service.erb
[Unit]
Description=<%= @description %>
After=syslog.target

[Service]
User=<%= @user %>
ExecStart=/usr/bin/java <%= @java_opts %> -jar <%= @jar_path %> --server.port=<%= @port %> --logging.path=<%= @logging_directory %> <%= @boot_opts %>
SuccessExitStatus=143
WorkingDirectory=/opt/spring-boot/<%= @new_resource_name %>

[Install]
WantedBy=multi-user.target
